# hw4cx/etc/ppc860-unified/BIVME2/rc.local
# Note:
#     This file is designed to be called from mounthost.sh

echo "Starting rrund"
$HOMEDIR/rrund -p 8002 $HOMEDIR/bivme2_ .drvlet -d >/dev/null 2>/dev/null &
